# templates/default/etcd.conf.erb
description "etcd service registry"

start on started networking
stop on shutdown

exec etcd -name 'etcd0' -data-dir '/tmp/etcd' \
                -advertise-client-urls http://`hostname -i`:2379,http://`hostname -i`:4001 \
                -listen-client-urls http://0.0.0.0:2379,http://0.0.0.0:4001 \
                -initial-advertise-peer-urls http://`hostname -i`:2380 \
                -listen-peer-urls http://0.0.0.0:2380 \
                -initial-cluster-token etcd-cluster-1 \
                -initial-cluster etcd0=http://`hostname -i`:2380 \
                -initial-cluster-state new
